# discoverd/etc/inittab
console::respawn:-/bin/sh
::ctrlaltdel:/sbin/reboot
